# init.d/health-agent.Debian-7
#! /bin/bash --posix
#
# health-agent		Start up the health-agent daemon
#
# chkconfig: 2345 95 25
# description: health-agent publishes system health metrics.

### BEGIN INIT INFO
# Provides:		health-agent
# Required-Start:	$local_fs $network $syslog
# Required-Stop:	$local_fs $network $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	System health agent
### END INIT INFO

# /etc/init.d/health-agent: start and stop the system health daemon

# CentOS-6 doesn't provide this, so provide a default.
log_daemon_msg ()
{
    echo "$@"
}

. /lib/lsb/init-functions

umask 022

readonly default_log_dir='/var/log/health-agent'
readonly LOOP_PIDFILE='/var/run/health-agent.loop.pid'
readonly PIDFILE='/var/run/health-agent.pid'

DAEMON='/usr/local/sbin/health-agent'
LOG_DIR="$default_log_dir"
LOG_QUOTA=
LOGBUF_LINES=
PROBE_INTERVAL=
USERNAME=

PROG_ARGS=

[ -f /etc/default/health-agent ] && . /etc/default/health-agent

test -x "$DAEMON" || exit 0

export PATH="${PATH:+$PATH:}/usr/local/bin:/usr/local/sbin:/usr/sbin:/sbin"

mkdir -m 0750 -p "$LOG_DIR"

if [ -n "$LOG_DIR" ] && [ "$LOG_DIR" != "$default_log_dir" ]; then
    PROG_ARGS="$PROG_ARGS -logDir=$LOG_DIR"
fi

if [ -n "$LOG_QUOTA" ]; then
    PROG_ARGS="$PROG_ARGS -logQuota=$LOG_QUOTA"
fi

if [ -n "$LOGBUF_LINES" ]; then
    PROG_ARGS="$PROG_ARGS -logbufLines=$LOGBUF_LINES"
fi

if [ -n "$PROBE_INTERVAL" ]; then
    PROG_ARGS="$PROG_ARGS -probeInterval=$PROBE_INTERVAL"
fi

do_start ()
{
    if [ -z "$USERNAME" ]; then
	"$DAEMON" $PROG_ARGS
	return
    fi
    chown "$USERNAME" "$LOG_DIR"
    start-stop-daemon --start --quiet --pidfile "$PIDFILE" \
		      --exec "$DAEMON" --chuid "$USERNAME" --make-pidfile \
		      --nicelevel 1 -- \
		      $PROG_ARGS
}

start_loop ()
{
    echo "$BASHPID" > "$LOOP_PIDFILE"
    while true; do
	do_start
	rm -f "$PIDFILE"
	sleep 1
    done
}

case "$1" in
  start)
	log_daemon_msg "Starting system health agent" "health-agent" || true
	(start_loop < /dev/null &> /dev/null &)
	;;
  stop)
	log_daemon_msg "Stopping system health agent" "health-agent" || true
	if [ -s "$LOOP_PIDFILE" ]; then
	     kill -KILL $(cat "$LOOP_PIDFILE")
	     rm -f "$LOOP_PIDFILE"
	fi
	if [ -s "$PIDFILE" ]; then
	    kill -TERM $(cat "$PIDFILE")
	    rm -f "$PIDFILE"
	else
	    killall -o 1s health-agent
	fi
	;;

  reload|force-reload)
	kill -HUP $(cat "$PIDFILE")
	;;

  restart)
	kill -HUP $(cat "$PIDFILE")
	;;

  start-if-down)
	[ -s "$LOOP_PIDFILE" ] && kill -0 $(cat "$LOOP_PIDFILE") &> /dev/null && exit
	killall -o 1s health-agent || true
	log_daemon_msg "Starting system health agent" "health-agent" || true
	(start_loop < /dev/null &> /dev/null &)
	;;

  *)
	log_action_msg "Usage: /etc/init.d/health-agent {start|stop|reload|force-reload|restart}" || true
	exit 1
esac

exit 0
